# build/src/debian/postinst
#!/bin/sh
#
# see: dh_installdeb(1)

#set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

USER=wekan
GROUP=wekan
WORKING_DIR=/opt/bin/wekan

check_user()
{
    if [ "$(id -u $USER 2>/dev/null)" = "" ] ; then
        useradd --user-group --system --home-dir $WORKING_DIR $USER
        chown -R $USER:$USER $WORKING_DIR
    fi
}

case "$1" in

    configure)
        check_user
    ;;

    upgrade)
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;

esac

if [ -x "/etc/init.d/wekan-oft-0" ]; then
    /etc/init.d/wekan-oft-0 stop
    update-rc.d -f wekan-oft-0 remove
    update-rc.d wekan-oft-0 defaults
    /etc/init.d/wekan-oft-0 start
fi

exit 0